# final/parts.d/appimged/profile.d/appimaged-start.sh
systemctl --user start appimaged.service || true
