# final/parts.d/appimged/profile.d/appimaged-start.sh
if puavo.profiles.list| grep -q "laptop,personal" ;then
   systemctl --user start appimaged.service || true
fi
