# final/parts.d/avahi/tree/etc/systemd/system/puavo-avahi-publish.service
[Unit]
Description=Announce Puavo Devices
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/puavo-avahi-publish
[Install]
WantedBy=multi-user.target
